# packaging/debian/terminusd.postrm
#!/bin/sh
set -e

case "$1" in
    remove|purge)
        if command -v systemctl >/dev/null 2>&1; then
            systemctl disable --now terminusd.service >/dev/null 2>&1 || true
        fi
        rm -f /etc/terminus.d/notifyusers.conf
        rm -f /etc/terminus.d/applyupdates.conf
        rm -f /etc/terminus.d/terminateusersessions.conf
        rm -f /usr/libexec/terminusd/examples/example-shutdown-notify.sh
        rm -f /usr/libexec/terminusd/examples/example-package-updates.sh
        rm -f /usr/libexec/terminusd/examples/example-persistent-shutdown-guard.sh
        rm -f /usr/libexec/terminusd/examples/example-oneshot-shutdown-guard.sh
        rm -f /usr/libexec/terminusd/examples/example-terminate-user-sessions.sh
        rmdir /etc/terminus.d >/dev/null 2>&1 || true
        rmdir /usr/libexec/terminusd/examples >/dev/null 2>&1 || true
        rmdir /usr/libexec/terminusd >/dev/null 2>&1 || true
        ;;
esac

#DEBHELPER#

exit 0
