# daemon/system-update-inhibitor.service
[Unit]
Description=System Update Inhibitor Daemon
After=network.target

[Service]
Type=forking
ExecStart=/usr/sbin/system-update-inhibitor
Restart=on-failure

[Install]
WantedBy=multi-user.target
